# xdg/upstart/a2jmidid.conf
description "ALSA to JACK MIDI Bridge"

start on started jackdbus
stop on stopping jackdbus

respawn
exec a2jmidid dbus
