# xdg/upstart/brillo.conf
description "brillo: persist brightness across logins"

start on session-up
stop on session-end

pre-start exec brillo -I
post-stop exec brillo -O
