# xdg/upstart/jackdbus.conf
description "JACK Audio Connection Kit D-Bus Server"

start on studio-up
stop on studio-end

respawn
exec jackdbus auto
