# xdg/upstart/jmcore.conf
start on started jackdbus
stop on stopping jackdbus

respawn
exec jmcore
