# xdg/upstart/ladiconfd.conf
description "LADI Configuration Daemon"

start on started jackdbus
stop on stopping jackdbus

exec ladiconfd
