# xdg/upstart/ladishd.conf
start on started ladiconfd
stop on stopping ladiconfd

respawn
exec ladishd
