# xdg/upstart/redshift.conf
start on started sway
stop on session-end

# unfortunately, sway and redshift race to grab the drm device
pre-start exec sleep 1
respawn
exec redshift
