# xdg/upstart/swaylock.conf
description "swaylock: lock the session before sleep"
start on dbus SIGNAL='PrepareForSleep' BUS='system' INTERFACE='org.freedesktop.login1.Manager'
exec swaylock
