# xdg/upstart/swaylock.conf
description "swaylock: lock the session before sleep"
start on started sway or stopping woke
respawn
normal exit 0
exec swaylock -i /usr/share/backgrounds/default
